# .config/systemd/user/hibernate.target
[Unit]
Description=Hibernate
StopWhenUnneeded=yes