# .config/systemd/user/hybrid-sleep.target
[Unit]
Description=Hybrid Sleep
StopWhenUnneeded=yes