# .config/systemd/user/suspend-then-hibernate.target
[Unit]
Description=Suspend Then Hibernate
StopWhenUnneeded=yes